# S-reglas-udev/buskill.rules
ACTION=="add", SUBSYSTEM=="block", RUN+="/root/insertar.sh" 
ACTION=="remove", SUBSYSTEM=="block", RUN+="/root/retirar.sh" 
