# S-reglas-udev/buskill.rules
ACTION=="remove", SUBSYSTEM=="block", RUN+="/root/retirar.sh" 
